# deploy/paasd.service
[Unit]
Description=paasd - Agentic PaaS daemon
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=paasd
Group=paasd
SupplementaryGroups=docker
EnvironmentFile=-/etc/default/paasd
ExecStart=/usr/local/bin/paasd serve \
  --port 8080 \
  --db-path /var/lib/paasd/paasd.db \
  --master-key-path /var/lib/paasd/master.key
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=paasd

StateDirectory=paasd

# Security hardening
UMask=0077
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes
ProtectSystem=full
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
PrivateDevices=yes
RestrictNamespaces=yes
ProtectKernelLogs=yes
ReadWritePaths=/var/lib/paasd /etc/traefik/dynamic/paasd.yml

[Install]
WantedBy=multi-user.target
